# macos/soundspace-decorum/crontab
* * * * * volume-auto.sh > volume.log 2>$1
